# standalone/systemd/mass-doc-to-pdf-api.service.in
[Unit]
Description=mass-doc-to-pdf API
After=network.target

[Service]
Type=simple
WorkingDirectory=__ROOT__/apps/api
EnvironmentFile=__ROOT__/.env.standalone
# Run the prebuilt server directly (build.sh + init-db.sh handle generate/migrate).
ExecStart=/usr/bin/node dist/server.js
Restart=always
RestartSec=5
User=__USER__
Group=__USER__

[Install]
WantedBy=multi-user.target
